# mission-control/mission-control.service
[Unit]
Description=Mission Control - Second Brain Dashboard
After=network.target

[Service]
Type=simple
User=ehsanjso
WorkingDirectory=/home/ehsanjso/clawd/mission-control/.next/standalone
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

# Memory limit (good for Pi)
MemoryMax=512M

[Install]
WantedBy=multi-user.target
